# wayland-screenlock-proxy.service.in
[Unit]
Description=Screen locker proxy service
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
ExecStart=@libexec@/lwrap
Type=exec

[Install]
WantedBy=graphical-session.target
