# wayland-screenlock-proxy.service.in
[Unit]
Description=Screen locker proxy service
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
ExecStart=@libexec@/wayland-screenlock-proxy
BusName=org.example.wayland_screenlock_proxy

[Install]
WantedBy=graphical-session.target
